# output/linux/service/opencvr-ubuntu.conf
author "xsmart@veyesys.cpm"
description "start and stop OpenCVR Stor for Ubuntu (upstart)"
version "1.0"

start on started networking
stop on runlevel [!2345]

env APPUSER="user"
env APPDIR="/home/user/work/opencvr/OpenCVR-Ubuntu-14.04-64bit-r2.0.1-alpha3/bin"

env LD_LIBRARY_PATH=$APPDIR/../lib/:$LD_LIBRARY_PATH
env QT_QPA_FONTDIR=$APPDIR/../lib/fonts/
env QT_QPA_PLATFORM_PLUGIN_PATH=$APPDIR/../plugins/platforms/

env APPBIN="OpenCVRStor"
env APPARGS="-platform offscreen "

respawn

script
  exec su - $APPUSER -c "$APPDIR/$APPBIN $APPARGS"
end script
