# output/linux/service/rapidvms-ubuntu.conf
author "support@veyesys.cpm"
description "start and stop Rapidvms Stor for Ubuntu (upstart)"
version "1.0"

start on started networking
stop on runlevel [!2345]

env APPUSER="user"
env APPDIR="/home/user/work/rapidvms/Rapidvms-Ubuntu-14.04-64bit-r2.0.1-alpha3/bin"

env LD_LIBRARY_PATH=$APPDIR/../lib/:$LD_LIBRARY_PATH

env APPBIN="RapidvmsStor"
env APPARGS="-platform offscreen "

respawn

script
  exec su - $APPUSER -c "$APPDIR/$APPBIN $APPARGS"
end script
